# contrib/archlinux/ebusd.install
post_install() {
  echo "Instructions:"
  echo "1. Edit /etc/conf.d/ebusd if necessary"
  echo "   (especially if your device is not /dev/ttyACM0)"
  echo "2. Start the daemon with 'systemctl start ebusd'"
  echo "3. Check the log file /var/log/ebusd.log"
  echo "4. Make the daemon autostart with 'systemctl enable ebusd'"
}

