# contrib/gentoo/init.d/ebusd
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

command="/usr/bin/ebusd"
command_args="${EBUSD_OPTS}"
start_stop_daemon_args="--quiet"
description="ebusd, the daemon for communication with eBUS heating systems."

depend() {
	need localmount
	use logger
}
